# systemd/term_planner-notify.service
[Unit]
Description=TermPlanner Notification Service

[Service]
Type=oneshot
ExecStart=%h/.local/bin/term_planner --notify
