# systemd/term_planner-notify.service
[Unit]
Description=TermPlanner Notification Service

[Service]
Type=oneshot
ExecStart=/usr/bin/term_planner --notify

